# certbot-renewal-test.service
[Unit]
Description=Certbot Renewal Test
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/certbot-renewal.sh --dry-run
User=root
